# services/spotict.service
[Unit]
Description=spotict display cover art service
After=graphical.target

[Service]
X11InitCheck=true
Type=simple
ExecStartPre=/bin/sleep 10
SyslogIdentifier=spotict

[Install]
WantedBy=multi-user.target
