# services/spotict.service
[Unit]
Description=spotict display cover art service
After=display-manager.service
Wants=display-manager.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 60
SyslogIdentifier=spotict

[Install]
WantedBy=multi-user.target
