# services/spotict.service
[Unit]
Description=spotict display cover art service

[Service]
Type=simple
X11InitCheck=true
ExecStartPre=/bin/sleep 60
SyslogIdentifier=spotict

[Install]
WantedBy=multi-user.target
