# services/spotict.service
[Unit]
Description=spotict display cover art service

[Service]
User=will
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=spotipi
Environment="DISPLAY=:0"
ExecStartPre=/bin/sleep 10
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
